# install/savecraft.service
# Reference only — the installer generates the actual unit file inline
# with paths derived from SAVECRAFT_APP_NAME (see install.sh).
[Unit]
Description=Savecraft Daemon
After=network-online.target

[Service]
ExecStart=%h/.local/bin/savecraft-daemon
Restart=always
RestartSec=5
EnvironmentFile=-%h/.config/savecraft/env

ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/savecraft %h/.cache/savecraft %h/.local/bin

NoNewPrivileges=yes
PrivateTmp=yes

RestrictAddressFamilies=AF_INET AF_INET6

[Install]
WantedBy=default.target
